# systemd/nightsnack.timer
[Unit]
Description=Run nightsnack every 5min

[Timer]
OnBootSec=5min
OnUnitActiveSec=5min

[Install]
WantedBy=timers.target
